# android-mtp@.service
# Android MTP mount service

[Service]
Type=forking
ExecStart=/usr/bin/android-mtp mount %I
ExecStop=/usr/bin/android-mtp unmount %I
